# systemd/org.deepin.dde.Launcher1.service.in
# SPDX-FileCopyrightText: 2023 UnionTech Software Technology Co., Ltd.
#
# SPDX-License-Identifier: GPL-3.0-or-later

[Unit]
Description=DDE launchpad background service

Requisite=dde-session-pre.target
After=dde-session-pre.target

Requisite=dde-session-initialized.target
PartOf=dde-session-initialized.target
Before=dde-session-initialized.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/dde-launchpad
Type=dbus
BusName=org.deepin.dde.Launcher1
